# projects/ops/com.paperclip.daily-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.paperclip.daily-backup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/fabiodomingues/miniconda3/bin/python3</string>
        <string>/Users/fabiodomingues/Desktop/Projects/paperclip/projects/ops/scripts/run-daily-backup.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/fabiodomingues</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/fabiodomingues/.paperclip/logs/daily-backup-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/fabiodomingues/.paperclip/logs/daily-backup-stderr.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
