# projects/whatsapp/com.paperclip.whatsapp-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.paperclip.whatsapp-report</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/osascript</string>
        <string>-e</string>
        <string>do shell script "export PATH=/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin; /usr/bin/python3 /Users/fabiodomingues/Desktop/Projects/paperclip/projects/whatsapp/campaign_report.py >> /Users/fabiodomingues/.paperclip/logs/whatsapp-report-stdout.log 2>&amp;1"</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/fabiodomingues/.paperclip/logs/whatsapp-report-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/fabiodomingues/.paperclip/logs/whatsapp-report-stderr.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
